# scheduler/ru.secondbrain.digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  ru.secondbrain.digest.plist — LaunchAgent проактивного слоя «Второго мозга».

  ШАБЛОН: замени плейсхолдеры __PUBLIC_REPO__ (и при необходимости пути логов)
  на абсолютные пути своей машины ПЕРЕД установкой. Инструкции по установке —
  в setup/SETUP.md (launchctl bootstrap gui/$(id -u) ...). НЕ запускай отсюда
  автоматически — это человеческий шаг.

  ПОЧЕМУ ИМЕННО ТАК (research/proactive-scheduling.md + ADR-0007 §2):
    • LaunchAgent (этот файл кладётся в ~/Library/LaunchAgents), НЕ LaunchDaemon:
      работает как залогиненный GUI-юзер → есть сеть и доступ к `claude login`-сессии
      / keychain (ADR-0008/0009), sudo не нужен. LaunchDaemon (root) — для будущего
      always-on Mac Mini (либо remote Claude routines — см. README.md).
    • ОДИН plist с ДВУМЯ триггерами:
        - StartInterval 1800с (= каждые 30 мин): частый идемпотентный sweep, чтобы
          пробуждение из сна рядом с due-временем триггерило проверку. launchd при
          пробуждении запускает пропущенный интервал (cron — скипает) и КОАЛЕСЦИРУЕТ
          несколько пропущенных в одно событие — поэтому job делает один sweep по
          ВСЕМ due, а не таймер-на-напоминание.
        - StartCalendarInterval (Hour 8, Minute 0): фиксированный утренний дайджест.
    • RunAtLoad = true: стрелять и на login/load (покрывает gap после загрузки,
      который wake-catch-up не покрывает).
    • Дедуп от двойного запуска (коалесцирование) — в коде (status/last_fired),
      не в launchd.

  ОГРАНИЧЕНИЕ (осознанное, ADR-0005): слот, пропущенный при ПОЛНОСТЬЮ ВЫКЛЮЧЕННОМ
  Mac, не отыгрывается на загрузке (ловится только сон). Настоящий always-on —
  фаза Mac Mini / LaunchDaemon.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Уникальный лейбл агента (совпадает с именем файла без .plist). -->
    <key>Label</key>
    <string>ru.secondbrain.digest</string>

    <!-- Что запускать: обёртка run_sweep.sh (грузит .env, caffeinate, python). -->
    <!-- Замени __PUBLIC_REPO__ на абсолютный путь к репо personal-llm-wiki. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PUBLIC_REPO__/scheduler/run_sweep.sh</string>
    </array>

    <!-- Триггер 1: частый sweep каждые 5 минут — due-напоминания доставляются скоро после
         срока (предчек дешёвый, без движка → claude -p спавнится ТОЛЬКО если что-то due,
         совместимо с ADR-0009). Тюнингуй под баланс отзывчивость/батарея: 1800=батарея,
         300=отзывчивые напоминания. -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Триггер 2: фиксированный утренний дайджест в 08:00 локального времени. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Стрелять также при загрузке/логине агента (покрывает post-boot gap). -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Не перезапускать мгновенно при выходе: sweep одноразовый, не демон. -->
    <key>KeepAlive</key>
    <false/>

    <!-- Рабочая директория = корень публичного репо (на всякий случай). -->
    <key>WorkingDirectory</key>
    <string>__PUBLIC_REPO__</string>

    <!-- Логи stdout/stderr (structlog/предчек). Поправь пути при необходимости. -->
    <key>StandardOutPath</key>
    <string>__PUBLIC_REPO__/scheduler/digest.out.log</string>
    <key>StandardErrorPath</key>
    <string>__PUBLIC_REPO__/scheduler/digest.err.log</string>

    <!-- Минимальное окружение. CONTENT_ROOT/секреты подтянет run_sweep.sh из
         приватного .env; здесь только PATH, чтобы нашлись bash/python/caffeinate.
         При нестандартном питоне допиши его bin-каталог в PATH. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <!-- Не давать процессу зависнуть навечно: жёсткий таймаут (сек). launchd
         пошлёт SIGTERM, если sweep не завершился (таймаут `claude -p` страхуется
         и в коде моста — ADR-0008). -->
    <key>ExitTimeOut</key>
    <integer>300</integer>

    <!-- Низкий приоритет (фоновая задача, не мешать интерактиву). -->
    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
