# bridge/ru.secondbrain.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ru.secondbrain.bridge.plist — LaunchAgent для моста Telegram <-> движок (Claude).

  Держит node-процесс моста живым (KeepAlive: автоперезапуск при падении) и
  стартует при логине. LaunchAgent (per-user), а НЕ LaunchDaemon — в v1 хост это
  MacBook владельца, без always-on-root-демона ([ADR-0005], [ADR-0007]).

  TypeScript-порт ([ADR-0012]): мост — Fastify на Node (dist/bridge/main.js), не
  FastAPI/uvicorn. СНАЧАЛА собери: pnpm install && pnpm build (tsc → dist/).

  УСТАНОВКА (детали — в setup/SETUP.md, ставит человек, не код):
    1. Замени плейсхолдеры __...__ на реальные абсолютные пути этой машины.
    2. Скопируй в ~/Library/LaunchAgents/ru.secondbrain.bridge.plist
    3. launchctl load -w ~/Library/LaunchAgents/ru.secondbrain.bridge.plist
    4. Проверка: launchctl list | grep ru.secondbrain.bridge
                 curl http://127.0.0.1:8080/health

  ВАЖНО про секреты: НЕ хардкодим токены в plist. Мост читает их из .env в
  WorkingDirectory (dotenv-flow). Этот файл — публичный пример.
-->
<plist version="1.0">
<dict>
    <!-- Уникальная метка сервиса (reverse-DNS). -->
    <key>Label</key>
    <string>ru.secondbrain.bridge</string>

    <!-- Команда запуска: node поднимает Fastify-мост (dist/bridge/main.js).
         main.ts слушает 127.0.0.1:BRIDGE_PORT (дефолт 8080) — перед Cloudflare
         Tunnel мост слушает только loopback. Замени __NODE_BIN__ на путь к node
         (напр. /opt/homebrew/bin/node), __REPO_DIR__ — на корень публичного репо. -->
    <key>ProgramArguments</key>
    <array>
        <string>__NODE_BIN__</string>
        <string>__REPO_DIR__/dist/bridge/main.js</string>
    </array>

    <!-- Рабочая директория = корень публичного репо (тут dist/ и .env, gitignored).
         main.ts грузит .env через dotenv-flow из cwd. -->
    <key>WorkingDirectory</key>
    <string>__REPO_DIR__</string>

    <!-- Минимальное окружение. Полную конфигурацию мост берёт из .env;
         PATH нужен, чтобы дочерний официальный бинарь `claude` нашёлся (поправь под
         Homebrew-префикс: Apple Silicon — /opt/homebrew/bin, Intel — /usr/local/bin). -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Стартовать при логине пользователя. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Автоперезапуск при падении. -->
    <key>KeepAlive</key>
    <true/>

    <!-- Не душить рестартами: не чаще раза в 10с. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Логи stdout/stderr в файлы (pino пишет в stdout). Создай папку logs/. -->
    <key>StandardOutPath</key>
    <string>__REPO_DIR__/logs/bridge.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_DIR__/logs/bridge.err.log</string>
</dict>
</plist>
