# bridge/ru.secondbrain.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ru.secondbrain.bridge.plist — LaunchAgent для моста Telegram <-> движок (Claude).

  Держит uvicorn-процесс моста живым (KeepAlive: автоперезапуск при падении) и
  стартует при логине. Это LaunchAgent (per-user), а НЕ LaunchDaemon — в v1 хост
  это MacBook владельца, без always-on-root-демона ([ADR-0005], [ADR-0007]).

  NB про слои исполнения: этот plist обслуживает РЕАКТИВНЫЙ слой (Telegram → bridge →
  движок). Плановый слой (compile/digest/lint/web-research) — отдельно: локально
  через scheduler/ (launchd), а 24/7-апгрейд — remote Claude routines над приватным
  GitHub-репо (работают даже при спящем Mac) — [ADR-0005], [CONTEXT §4].

  УСТАНОВКА (детали — в setup/SETUP.md, ставит человек, не код):
    1. Замени ВСЕ плейсхолдеры __...__ на реальные абсолютные пути этой машины.
    2. Скопируй файл в ~/Library/LaunchAgents/ru.secondbrain.bridge.plist
    3. launchctl load -w ~/Library/LaunchAgents/ru.secondbrain.bridge.plist
    4. Проверка: launchctl list | grep ru.secondbrain.bridge
                 curl http://127.0.0.1:8080/health

  ВАЖНО про секреты: НЕ хардкодим токены в этот plist. Мост читает их из .env в
  WorkingDirectory (или из EnvironmentVariables, если решишь так). Этот файл —
  публичный пример, реальные значения сюда не попадают.
-->
<plist version="1.0">
<dict>
    <!-- Уникальная метка сервиса (reverse-DNS). -->
    <key>Label</key>
    <string>ru.secondbrain.bridge</string>

    <!-- Команда запуска: uvicorn поднимает app:app на localhost:порт.
         Перед Cloudflare Tunnel мост слушает только loopback (127.0.0.1).
         Замени __PYTHON_BIN__ на путь к python из venv моста (см. SETUP.md),
         например /Users/example/llm-wiki/bridge/.venv/bin/python -->
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_BIN__</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>app:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8080</string>
    </array>

    <!-- Рабочая директория = папка bridge/ (тут лежат app.py и .env). -->
    <key>WorkingDirectory</key>
    <string>__BRIDGE_DIR__</string>

    <!-- Минимальное окружение. Полную конфигурацию мост берёт из .env в
         WorkingDirectory; здесь задаём только то, без чего uvicorn не стартует.
         PATH нужен, чтобы дочерний официальный бинарь `claude` нашёлся (поправь под
         Homebrew-префикс: Apple Silicon — /opt/homebrew/bin, Intel — /usr/local/bin;
         если claude установлен иначе — добавь его каталог в PATH). -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Стартовать при логине пользователя. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Автоперезапуск при падении (фаза-7 hardening pachca-плана). -->
    <key>KeepAlive</key>
    <true/>

    <!-- Не душить рестартами: не чаще раза в 10с. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Логи stdout/stderr в файлы (structlog пишет в stdout). Создай папку logs/. -->
    <key>StandardOutPath</key>
    <string>__BRIDGE_DIR__/logs/bridge.out.log</string>
    <key>StandardErrorPath</key>
    <string>__BRIDGE_DIR__/logs/bridge.err.log</string>
</dict>
</plist>
